# auto-report.cron
30 7/24 * * * /usr/bin/python /daily-report/app.py $XMU_USERNAME $XMU_PASSWORD check