# auto-report.cron
30 */24 * * * /usr/bin/python /path/to/checkin.py [username] [password]